# systemd/os-ialocal.service
[Unit]
Description=OS IA Local — Chat UI para Ollama
After=network.target ollama.service

[Service]
Type=simple
User=osserver
WorkingDirectory=/home/osserver/Proyectos_Antigravity/Integraciones_OS/ialocal
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
